# john-inspiron-7420/etc/init/click-system-hooks.conf
description "Run Click system-level hooks"
author "Colin Watson <cjwatson@ubuntu.com>"

start on filesystem

task

exec /usr/bin/click hook run-system
